# vsetkonajlepsie.service
[Unit]
Description=Happy birthday, sucker!

[Service]
ExecStart=/usr/bin/vsetkonajlepsie.sh

[Install]
WantedBy=multi-user.target 
